# etc/xdg/autostart/system-kernel.desktop
[Desktop Entry]
Type=Application
Terminal=false
Name=System Old Kernel
Exec=/usr/bin/system-kernel
